# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

; telegram-bot-api starts first (lower priority number = earlier) and
; must stay up for startsecs before supervisord considers it "started"
; and moves on to the next process. This is a best-effort head start,
; not a true readiness check (supervisord doesn't probe the HTTP port)
; — the real safety net is the retry/backoff loop in main.py, which
; tolerates this process not being ready yet regardless.
[program:telegram-bot-api]
command=sh -c 'telegram-bot-api --local -d /app/telegram-bot-api-data --http-port 8081 --api-id="${TELEGRAM_API_ID}" --api-hash="${TELEGRAM_API_HASH}"'
autostart=true
autorestart=true
startsecs=3
startretries=5
priority=1
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
user=root

[program:fastapi-proxy]
command=python3 main.py
autostart=true
autorestart=true
startsecs=3
startretries=5
priority=10
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
user=root
